# ansible/roles/nextcloud_corupt_files_reporter/files/nc_reports.timer
[Unit]
Description=Nextcloud Reports timer

[Timer]
OnCalendar=01:00:00
Unit=nc_reports.service

[Install]
WantedBy=multi-user.target
