# ansible/roles/nextcloud_corupt_files_reporter/files/nc_reports.service
[Unit]
Description=Run Nextcloud Corrupt Files Reporting Script

[Service]
ExecStart=/usr/bin/bash -c '/etc/nc_reports.sh'
