# ansible/roles/syncoid_backup/files/backup.service
[Unit]
Description=Run Syncoid Backup

[Service]
ExecStart=/usr/bin/bash -c '/etc/backup.sh'
